# systemd/pipeline.timer
[Unit]
Description=Requirements Creator パイプライン 毎日JST実行

[Timer]
OnCalendar=*-*-* 02:00:00
OnCalendar=*-*-* 02:30:00
OnCalendar=*-*-* 03:00:00
OnCalendar=*-*-* 03:30:00
OnCalendar=*-*-* 04:00:00
OnCalendar=*-*-* 05:00:00
TimezoneOfTimer=Asia/Tokyo
Persistent=true

[Install]
WantedBy=timers.target
