# systemd/self-healing.service
[Unit]
Description=Requirements Creator 自己修復実行
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/develop/requirements_creator
ExecStart=%h/develop/requirements_creator/scripts/self-healing-run.sh
Environment=PATH=%h/.local/share/pnpm:%h/.nvm/versions/node/v22.17.0/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
