# .devcontainer/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:viewer]
command=pnpm viewer
directory=/workspace
autostart=true
autorestart=true
stdout_logfile=/workspace/logs/viewer-stdout.log
stderr_logfile=/workspace/logs/viewer-stderr.log
